# PS Login Monitor - RDP.xml
<?xml version="1.0" encoding="UTF-16"?>
<Task version="1.2" xmlns="http://schemas.microsoft.com/windows/2004/02/mit/task">
  <RegistrationInfo>
    <Description>Triggered by event 4625 in the Windows event log which indicate failed Windows login.</Description>
    <URI>\PS Login Monitor - RDP</URI>
  </RegistrationInfo>
  <Triggers>
    <EventTrigger>
      <Enabled>true</Enabled>
      <Subscription>&lt;QueryList&gt;&lt;Query Id="0" Path="Security"&gt;&lt;Select Path="Security"&gt;*[System[Provider[@Name='Microsoft-Windows-Security-Auditing'] and EventID=4625]]&lt;/Select&gt;&lt;/Query&gt;&lt;/QueryList&gt;</Subscription>
	  <!-- RDP event data. -->
      <ValueQueries>
        <Value name="EventID">Event/System/EventID</Value>
        <Value name="SubjectUserSid">Event/EventData/Data[@Name='SubjectUserSid']</Value>
        <Value name="SubjectUserName">Event/EventData/Data[@Name='SubjectUserName']</Value>
        <Value name="SubjectDomainName">Event/EventData/Data[@Name='SubjectDomainName']</Value>
        <Value name="SubjectLogonId">Event/EventData/Data[@Name='SubjectLogonId']</Value>
        <Value name="TargetUserSid">Event/EventData/Data[@Name='TargetUserSid']</Value>
        <Value name="TargetUserName">Event/EventData/Data[@Name='TargetUserName']</Value>
        <Value name="TargetDomainName">Event/EventData/Data[@Name='TargetDomainName']</Value>
        <Value name="Status">Event/EventData/Data[@Name='Status']</Value>
        <Value name="FailureReason">Event/EventData/Data[@Name='FailureReason']</Value>
        <Value name="SubStatus">Event/EventData/Data[@Name='SubStatus']</Value>
        <Value name="LogonType">Event/EventData/Data[@Name='LogonType']</Value>
        <Value name="LogonProcessName">Event/EventData/Data[@Name='LogonProcessName']</Value>
        <Value name="AuthenticationPackageName">Event/EventData/Data[@Name='AuthenticationPackageName']</Value>
        <Value name="WorkstationName">Event/EventData/Data[@Name='WorkstationName']</Value>
        <Value name="TransmittedServices">Event/EventData/Data[@Name='TransmittedServices']</Value>
        <Value name="LmPackageName">Event/EventData/Data[@Name='LmPackageName']</Value>
        <Value name="KeyLength">Event/EventData/Data[@Name='KeyLength']</Value>
        <Value name="ProcessId">Event/EventData/Data[@Name='ProcessId']</Value>
        <Value name="ProcessName">Event/EventData/Data[@Name='ProcessName']</Value>
        <Value name="IpAddress">Event/EventData/Data[@Name='IpAddress']</Value>
        <Value name="IpPort">Event/EventData/Data[@Name='IpPort']</Value>
      </ValueQueries>
    </EventTrigger>
  </Triggers>
  <Principals>
    <Principal id="Author">
      <UserId>S-1-5-18</UserId>
      <RunLevel>HighestAvailable</RunLevel>
    </Principal>
  </Principals>
  <Settings>
    <MultipleInstancesPolicy>Parallel</MultipleInstancesPolicy>
    <DisallowStartIfOnBatteries>false</DisallowStartIfOnBatteries>
    <StopIfGoingOnBatteries>false</StopIfGoingOnBatteries>
    <AllowHardTerminate>true</AllowHardTerminate>
    <StartWhenAvailable>false</StartWhenAvailable>
    <RunOnlyIfNetworkAvailable>false</RunOnlyIfNetworkAvailable>
    <IdleSettings>
      <StopOnIdleEnd>true</StopOnIdleEnd>
      <RestartOnIdle>false</RestartOnIdle>
    </IdleSettings>
    <AllowStartOnDemand>false</AllowStartOnDemand>
    <Enabled>true</Enabled>
    <Hidden>false</Hidden>
    <RunOnlyIfIdle>false</RunOnlyIfIdle>
    <WakeToRun>false</WakeToRun>
    <ExecutionTimeLimit>PT1H</ExecutionTimeLimit>
    <Priority>0</Priority>
  </Settings>
  <Actions Context="Author">
    <Exec>
      <Command>powershell.exe</Command>
      <Arguments>-command "&amp; {. .\PSLoginMonitor.ps1; On-FailedRdpLogin '$(EventID)' '$(SubjectUserSid)' '$(SubjectUserName)' '$(SubjectDomainName)' '$(SubjectLogonId)' '$(TargetUserSid)' '$(TargetUserName)' '$(TargetDomainName)' '$(Status)' '$(FailureReason)' '$(SubStatus)' '$(LogonType)' '$(LogonProcessName)' '$(AuthenticationPackageName)' '$(WorkstationName)' '$(TransmittedServices)' '$(LmPackageName)' '$(KeyLength)' '$(ProcessId)' '$(ProcessName)' '$(IpAddress)' '$(IpPort)'}"</Arguments>
	  <WorkingDirectory></WorkingDirectory>
    </Exec>
  </Actions>
</Task>
